# felicity/deploy-scripts/systemd/felicity-tasks.target
[Unit]
Description=Felicity Lims RQ worker pool

[Install]
WantedBy=multi-user.target